# cgminer/root-files/etc/uci-defaults/01-cgminer
#!/bin/sh

/etc/init.d/cgminer enable

# Setup the cron job
echo "*/5 *  *   *   *     /usr/bin/cgminer-monitor" | crontab -
/etc/init.d/cron enable
/etc/init.d/cron stop
/etc/init.d/cron start
