# cgminer/root-files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# Load hasher bitstream
[ -e /dev/xdevcfg ] && cat /root/hasher.bit > /dev/xdevcfg

# Close ipv6
sleep 3
sysctl -p

exit 0
